# scripts/debian/postinst
if ! getent passwd cloud-manager >/dev/null; then
    adduser --system --group --force-badname --quiet \
        --home /run/cloud-manager \
        --shell /bin/bash \
        cloud-manager
fi

if ! [ -d "/run/cloud-manager" ]; then
    mkdir -p "/run/cloud-manager"
    chown cloud-manager:cloud-manager "/run/cloud-manager"
fi
if ! [ -f "/opt/cloud-manager/keys/cloud-manager.key" ]; then
    mkdir -p '/opt/cloud-manager/keys/'
    openssl req -x509 -newkey rsa:4096 -subj "/CN=$(hostname -I | cut -d" " -f1 | xargs)" -addext "subjectAltName=IP:$(hostname -I | cut -d" " -f1 | xargs),IP:127.0.0.1,DNS:$(hostname)" -keyout /opt/cloud-manager/keys/cloud-manager.key -nodes -out /opt/cloud-manager/keys/cloud-manager.pem -sha256 -days 358000
fi
        
chmod -R 770 /opt/cloud-manager
chown -R cloud-manager:cloud-manager /opt/cloud-manager

if [ -f "/usr/lib/systemd/system/cloud-manager.service" ]; then
    rm -rf /usr/lib/systemd/system/cloud-manager.service  2>/dev/null || true
    systemctl disable cloud-manager.service 2>/dev/null || true
    systemctl stop cloud-manager.service 2>/dev/null || true
    systemctl daemon-reload 2>/dev/null || true
fi

echo """
[Unit]
Description=Fiber App Template (%I)
[Service]
Type=simple
WorkingDirectory=/opt/cloud-manager
ExecStart=/opt/cloud-manager/app -type=%i
Restart=always
RestartSec=10
SyslogIdentifier=cloud-manager
KillSignal=SIGINT
User=root
Group=root
[Install]
WantedBy=multi-user.target
    """ > /etc/systemd/system/cloud-manager@.service

systemctl daemon-reload
systemctl enable cloud-manager@admin.service cloud-manager@client.service
systemctl restart cloud-manager@admin.service cloud-manager@client.service